# supervisord.conf
[supervisord]
nodaemon=true
silent=true
logfile=/dev/stdout
pidfile=/var/run/supervisord.pid
user=root

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

# Socat SMTP listener (port 8025 -> 1025)
[program:socat-smtp]
command=socat TCP-LISTEN:8025,reuseaddr,fork TCP:127.0.0.1:1025
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
user=socat

# Socat IMAP listener (port 8143 -> 1143)
[program:socat-imap]
command=socat TCP-LISTEN:8143,reuseaddr,fork TCP:127.0.0.1:1143
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
user=socat

# ProtonMail Bridge
[program:protonmail-bridge]
command=/opt/protonmail/proton-bridge --cli
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
user=protonmail
environment=HOME=/home/protonmail
